# as4/80-penguin.rules
ENV{ID_INPUT_KEY}!="?", ENV{ID_MODEL}=="*Keyboard*", ACTION=="add", RUN+="/bin/device_added.sh"
