# as4/80-penguin.rules
SUBSYSTEMS=="input", ATTRS{name}=="* keyboard", ACTION=="add", RUN+="/bin/device_added.sh"
